# image/ports/contrib/wicd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/logrotate.d/
-rw-rw-r--	root/root	etc/logrotate.d/wicd.logrotate
drwxr-xr-x	root/root	etc/rc.d/
-rwxr--r--	root/root	etc/rc.d/wicd
drwxr-xr-x	root/root	etc/wicd/
-rw-rw-r--	root/root	etc/wicd/dhclient.conf.template.default
drwxr-xr-x	root/root	etc/wicd/encryption/
drwxr-xr-x	root/root	etc/wicd/encryption/templates/
-rw-rw-r--	root/root	etc/wicd/encryption/templates/active
-rw-rw-r--	root/root	etc/wicd/encryption/templates/active_wired
-rw-rw-r--	root/root	etc/wicd/encryption/templates/eap
-rw-rw-r--	root/root	etc/wicd/encryption/templates/eap-tls
-rw-rw-r--	root/root	etc/wicd/encryption/templates/leap
-rw-rw-r--	root/root	etc/wicd/encryption/templates/peap
-rw-rw-r--	root/root	etc/wicd/encryption/templates/peap-tkip
-rw-rw-r--	root/root	etc/wicd/encryption/templates/psu
-rw-rw-r--	root/root	etc/wicd/encryption/templates/ttls
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wep-hex
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wep-passphrase
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wep-shared
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wired_8021x
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wpa
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wpa-peap
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wpa-psk
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wpa-psk-hex
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wpa2-leap
-rw-rw-r--	root/root	etc/wicd/encryption/templates/wpa2-peap
drwxr-xr-x	root/root	etc/wicd/scripts/
-rw-rw-r--	root/root	etc/wicd/scripts/.empty_on_purpose
drwxr-xr-x	root/root	etc/wicd/scripts/postconnect/
-rw-rw-r--	root/root	etc/wicd/scripts/postconnect/.empty_on_purpose
drwxr-xr-x	root/root	etc/wicd/scripts/postdisconnect/
-rw-rw-r--	root/root	etc/wicd/scripts/postdisconnect/.empty_on_purpose
drwxr-xr-x	root/root	etc/wicd/scripts/preconnect/
-rw-rw-r--	root/root	etc/wicd/scripts/preconnect/.empty_on_purpose
drwxr-xr-x	root/root	etc/wicd/scripts/predisconnect/
-rw-rw-r--	root/root	etc/wicd/scripts/predisconnect/.empty_on_purpose
drwxr-xr-x	root/root	etc/xdg/
drwxr-xr-x	root/root	etc/xdg/autostart/
-rw-rw-r--	root/root	etc/xdg/autostart/wicd-tray.desktop
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/systemd/
drwxr-xr-x	root/root	lib/systemd/system/
-rw-rw-r--	root/root	lib/systemd/system/wicd.service
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxrwxr-x	root/root	usr/bin/wicd-cli
-rwxrwxr-x	root/root	usr/bin/wicd-client
-rwxrwxr-x	root/root	usr/bin/wicd-curses
-rwxrwxr-x	root/root	usr/bin/wicd-gtk
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/dbus-1/
drwxr-xr-x	root/root	usr/etc/dbus-1/system.d/
-rw-rw-r--	root/root	usr/etc/dbus-1/system.d/wicd.conf
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/python2.7/
drwxr-xr-x	root/root	usr/lib/python2.7/site-packages/
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd-1.7.4-py2.7.egg-info
drwxr-xr-x	root/root	usr/lib/python2.7/site-packages/wicd/
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/__init__.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/__init__.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/backend.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/backend.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/configmanager.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/configmanager.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/dbusmanager.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/dbusmanager.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/logfile.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/logfile.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/misc.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/misc.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/networking.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/networking.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/translations.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/translations.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/wnettools.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/wnettools.pyc
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/wpath.py
-rw-r--r--	root/root	usr/lib/python2.7/site-packages/wicd/wpath.pyc
drwxr-xr-x	root/root	usr/sbin/
-rwxrwxr-x	root/root	usr/sbin/wicd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/applications/
-rw-rw-r--	root/root	usr/share/applications/wicd.desktop
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-rw-r--	root/root	usr/share/dbus-1/system-services/org.wicd.daemon.service
drwxr-xr-x	root/root	usr/share/icons/
drwxr-xr-x	root/root	usr/share/icons/hicolor/
drwxr-xr-x	root/root	usr/share/icons/hicolor/128x128/
drwxr-xr-x	root/root	usr/share/icons/hicolor/128x128/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/128x128/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/16x16/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/192x192/
drwxr-xr-x	root/root	usr/share/icons/hicolor/192x192/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/192x192/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/22x22/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/24x24/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/32x32/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/36x36/
drwxr-xr-x	root/root	usr/share/icons/hicolor/36x36/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/36x36/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/48x48/
drwxr-xr-x	root/root	usr/share/icons/hicolor/48x48/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/48x48/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/64x64/
drwxr-xr-x	root/root	usr/share/icons/hicolor/64x64/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/64x64/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/72x72/
drwxr-xr-x	root/root	usr/share/icons/hicolor/72x72/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/72x72/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/96x96/
drwxr-xr-x	root/root	usr/share/icons/hicolor/96x96/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/96x96/apps/wicd-gtk.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/apps/
-rw-rw-r--	root/root	usr/share/icons/hicolor/scalable/apps/wicd-gtk.svg
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-rw-r--	root/root	usr/share/man/man1/wicd-client.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-rw-r--	root/root	usr/share/man/man5/wicd-manager-settings.conf.5.gz
-rw-rw-r--	root/root	usr/share/man/man5/wicd-wired-settings.conf.5.gz
-rw-rw-r--	root/root	usr/share/man/man5/wicd-wireless-settings.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-rw-r--	root/root	usr/share/man/man8/wicd-cli.8.gz
-rw-rw-r--	root/root	usr/share/man/man8/wicd-curses.8.gz
-rw-rw-r--	root/root	usr/share/man/man8/wicd.8.gz
drwxr-xr-x	root/root	usr/share/pixmaps/
-rw-rw-r--	root/root	usr/share/pixmaps/wicd-gtk.xpm
drwxr-xr-x	root/root	usr/share/wicd/
drwxr-xr-x	root/root	usr/share/wicd/backends/
-rw-rw-r--	root/root	usr/share/wicd/backends/be-external.py
-rw-r--r--	root/root	usr/share/wicd/backends/be-external.pyc
-rw-rw-r--	root/root	usr/share/wicd/backends/be-ioctl.py
-rw-r--r--	root/root	usr/share/wicd/backends/be-ioctl.pyc
drwxr-xr-x	root/root	usr/share/wicd/cli/
-rwxrwxr-x	root/root	usr/share/wicd/cli/wicd-cli.py
-rw-r--r--	root/root	usr/share/wicd/cli/wicd-cli.pyc
drwxr-xr-x	root/root	usr/share/wicd/curses/
-rwxrwxr-x	root/root	usr/share/wicd/curses/configscript_curses.py
-rw-r--r--	root/root	usr/share/wicd/curses/configscript_curses.pyc
-rw-rw-r--	root/root	usr/share/wicd/curses/curses_misc.py
-rw-r--r--	root/root	usr/share/wicd/curses/curses_misc.pyc
-rw-rw-r--	root/root	usr/share/wicd/curses/netentry_curses.py
-rw-r--r--	root/root	usr/share/wicd/curses/netentry_curses.pyc
-rw-rw-r--	root/root	usr/share/wicd/curses/prefs_curses.py
-rw-r--r--	root/root	usr/share/wicd/curses/prefs_curses.pyc
-rwxrwxr-x	root/root	usr/share/wicd/curses/wicd-curses.py
-rw-r--r--	root/root	usr/share/wicd/curses/wicd-curses.pyc
drwxr-xr-x	root/root	usr/share/wicd/daemon/
-rwxrwxr-x	root/root	usr/share/wicd/daemon/autoconnect.py
-rw-r--r--	root/root	usr/share/wicd/daemon/autoconnect.pyc
-rwxrwxr-x	root/root	usr/share/wicd/daemon/monitor.py
-rw-r--r--	root/root	usr/share/wicd/daemon/monitor.pyc
-rwxrwxr-x	root/root	usr/share/wicd/daemon/suspend.py
-rw-r--r--	root/root	usr/share/wicd/daemon/suspend.pyc
-rw-rw-r--	root/root	usr/share/wicd/daemon/wicd-daemon.py
-rw-r--r--	root/root	usr/share/wicd/daemon/wicd-daemon.pyc
drwxr-xr-x	root/root	usr/share/wicd/gtk/
-rwxrwxr-x	root/root	usr/share/wicd/gtk/configscript.py
-rw-r--r--	root/root	usr/share/wicd/gtk/configscript.pyc
-rw-rw-r--	root/root	usr/share/wicd/gtk/gui.py
-rw-r--r--	root/root	usr/share/wicd/gtk/gui.pyc
-rw-rw-r--	root/root	usr/share/wicd/gtk/guiutil.py
-rw-r--r--	root/root	usr/share/wicd/gtk/guiutil.pyc
-rw-rw-r--	root/root	usr/share/wicd/gtk/netentry.py
-rw-r--r--	root/root	usr/share/wicd/gtk/netentry.pyc
-rw-rw-r--	root/root	usr/share/wicd/gtk/prefs.py
-rw-r--r--	root/root	usr/share/wicd/gtk/prefs.pyc
-rw-rw-r--	root/root	usr/share/wicd/gtk/wicd-client.py
-rw-r--r--	root/root	usr/share/wicd/gtk/wicd-client.pyc
-rw-rw-r--	root/root	usr/share/wicd/gtk/wicd.ui
drwxr-xr-x	root/root	usr/share/wicd/icons/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/16x16/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/16x16/status/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/16x16/status/signal-100.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/16x16/status/signal-25.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/16x16/status/signal-50.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/16x16/status/signal-75.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/22x22/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/22x22/status/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/bad-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/bad-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-bad-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-bad-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-good-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-good-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-high-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-high-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-low-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/both-low-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/good-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/good-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/high-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/high-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-bad-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-bad-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-good-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-good-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-high-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-high-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-low-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/idle-low-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/low-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/low-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/no-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-bad-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-bad-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-good-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-good-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-high-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-high-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-low-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/receiving-low-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/signal-100.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/signal-25.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/signal-50.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/signal-75.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-bad-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-bad-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-good-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-good-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-high-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-high-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-low-signal-lock.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/transmitting-low-signal.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/22x22/status/wired.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/24x24/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/24x24/status/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/24x24/status/signal-100.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/24x24/status/signal-25.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/24x24/status/signal-50.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/24x24/status/signal-75.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/32x32/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/32x32/status/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/32x32/status/signal-100.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/32x32/status/signal-25.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/32x32/status/signal-50.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/32x32/status/signal-75.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/36x36/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/36x36/status/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/36x36/status/signal-100.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/36x36/status/signal-25.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/36x36/status/signal-50.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/36x36/status/signal-75.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/48x48/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/48x48/devices/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/48x48/devices/wired-gui.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/48x48/status/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/48x48/status/signal-100.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/48x48/status/signal-25.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/48x48/status/signal-50.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/48x48/status/signal-75.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/original/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/original/status/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/original/status/signal-100.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/original/status/signal-25.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/original/status/signal-50.png
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/original/status/signal-75.png
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/scalable/
drwxr-xr-x	root/root	usr/share/wicd/icons/hicolor/scalable/devices/
-rw-rw-r--	root/root	usr/share/wicd/icons/hicolor/scalable/devices/wired-gui.svg
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/wicd/
drwxr-xr-x	root/root	var/lib/wicd/configurations/
-rw-rw-r--	root/root	var/lib/wicd/configurations/.empty_on_purpose
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/wicd/
-rw-rw-r--	root/root	var/log/wicd/.empty_on_purpose
