# image/ports/opt/squid/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/squid
drwxr-xr-x	root/root	etc/squid/
-rw-r--r--	root/root	etc/squid/cachemgr.conf
-rw-r--r--	root/root	etc/squid/cachemgr.conf.default
-rw-r--r--	root/root	etc/squid/errorpage.css
-rw-r--r--	root/root	etc/squid/errorpage.css.default
-rw-r--r--	root/root	etc/squid/mime.conf
-rw-r--r--	root/root	etc/squid/mime.conf.default
-rw-r--r--	root/root	etc/squid/squid.conf
-rw-r--r--	root/root	etc/squid/squid.conf.default
-rw-r--r--	root/root	etc/squid/squid.conf.documented
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/purge
-rwxr-xr-x	root/root	usr/bin/squidclient
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/squid/
-rwxr-xr-x	root/root	usr/lib/squid/basic_db_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_fake_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_getpwnam_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_ncsa_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_nis_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_pam_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_pop3_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_radius_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_smb_auth
-rwxr-xr-x	root/root	usr/lib/squid/basic_smb_auth.sh
-rwxr-xr-x	root/root	usr/lib/squid/cachemgr.cgi
-rwxr-xr-x	root/root	usr/lib/squid/digest_file_auth
-rwxr-xr-x	root/root	usr/lib/squid/diskd
-rwxr-xr-x	root/root	usr/lib/squid/ext_delayer_acl
-rwxr-xr-x	root/root	usr/lib/squid/ext_file_userip_acl
-rwxr-xr-x	root/root	usr/lib/squid/ext_session_acl
-rwxr-xr-x	root/root	usr/lib/squid/ext_sql_session_acl
-rwxr-xr-x	root/root	usr/lib/squid/ext_time_quota_acl
-rwxr-xr-x	root/root	usr/lib/squid/ext_unix_group_acl
-rwxr-xr-x	root/root	usr/lib/squid/ext_wbinfo_group_acl
-rwxr-xr-x	root/root	usr/lib/squid/helper-mux
-rwxr-xr-x	root/root	usr/lib/squid/log_db_daemon
-rwxr-xr-x	root/root	usr/lib/squid/log_file_daemon
-rwxr-xr-x	root/root	usr/lib/squid/negotiate_wrapper_auth
-rwxr-xr-x	root/root	usr/lib/squid/ntlm_fake_auth
-rwxr-xr-x	root/root	usr/lib/squid/security_fake_certverify
-rwxr-xr-x	root/root	usr/lib/squid/storeid_file_rewrite
-rwxr-xr-x	root/root	usr/lib/squid/unlinkd
-rwxr-xr-x	root/root	usr/lib/squid/url_fake_rewrite
-rwxr-xr-x	root/root	usr/lib/squid/url_fake_rewrite.sh
-rwxr-xr-x	root/root	usr/lib/squid/url_lfs_rewrite
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/squid
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/purge.1.gz
-rw-r--r--	root/root	usr/share/man/man1/squidclient.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/basic_db_auth.8.gz
-rw-r--r--	root/root	usr/share/man/man8/basic_getpwnam_auth.8.gz
-rw-r--r--	root/root	usr/share/man/man8/basic_ncsa_auth.8.gz
-rw-r--r--	root/root	usr/share/man/man8/basic_pam_auth.8.gz
-rw-r--r--	root/root	usr/share/man/man8/basic_pop3_auth.8.gz
-rw-r--r--	root/root	usr/share/man/man8/basic_radius_auth.8.gz
-rw-r--r--	root/root	usr/share/man/man8/cachemgr.cgi.8.gz
-rw-r--r--	root/root	usr/share/man/man8/digest_file_auth.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ext_delayer_acl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ext_file_userip_acl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ext_session_acl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ext_sql_session_acl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ext_time_quota_acl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ext_unix_group_acl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ext_wbinfo_group_acl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/helper-mux.8.gz
-rw-r--r--	root/root	usr/share/man/man8/log_db_daemon.8.gz
-rw-r--r--	root/root	usr/share/man/man8/security_fake_certverify.8.gz
-rw-r--r--	root/root	usr/share/man/man8/squid.8.gz
-rw-r--r--	root/root	usr/share/man/man8/storeid_file_rewrite.8.gz
-rw-r--r--	root/root	usr/share/man/man8/url_lfs_rewrite.8.gz
drwxr-xr-x	root/root	usr/share/squid/
drwxr-xr-x	root/root	usr/share/squid/errors/
drwxr-xr-x	root/root	usr/share/squid/errors/templates/
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_ACCESS_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_ACL_TIME_QUOTA_EXCEEDED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_AGENT_CONFIGURE
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_AGENT_WPAD
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_CACHE_ACCESS_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_CACHE_MGR_ACCESS_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_CANNOT_FORWARD
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_CONFLICT_HOST
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_CONNECT_FAIL
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_DIR_LISTING
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_DNS_FAIL
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_ESI
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FORWARDING_DENIED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_DISABLED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_FAILURE
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_FORBIDDEN
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_NOT_FOUND
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_PUT_CREATED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_PUT_ERROR
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_PUT_MODIFIED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_FTP_UNAVAILABLE
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_GATEWAY_FAILURE
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_ICAP_FAILURE
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_INVALID_REQ
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_INVALID_RESP
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_INVALID_URL
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_LIFETIME_EXP
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_NO_RELAY
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_ONLY_IF_CACHED_MISS
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_PRECONDITION_FAILED
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_PROTOCOL_UNKNOWN
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_READ_ERROR
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_READ_TIMEOUT
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_SECURE_CONNECT_FAIL
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_SHUTTING_DOWN
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_SOCKET_FAILURE
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_TOO_BIG
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_UNSUP_HTTPVERSION
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_UNSUP_REQ
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_URN_RESOLVE
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_WRITE_ERROR
-rw-r--r--	root/root	usr/share/squid/errors/templates/ERR_ZERO_SIZE_OBJECT
-rw-r--r--	root/root	usr/share/squid/errors/templates/error-details.txt
drwxr-xr-x	root/root	usr/share/squid/icons/
-rw-r--r--	root/root	usr/share/squid/icons/SN.png
drwxr-xr-x	root/root	usr/share/squid/icons/silk/
-rw-r--r--	root/root	usr/share/squid/icons/silk/application.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/arrow_up.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/bomb.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/box.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/bricks.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/bullet_red.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/cd.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/chart_line.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/compress.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/computer_link.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/css.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/cup.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/database.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/database_table.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/drive_disk.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/film.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/film_key.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/folder.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/folder_table.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/image.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/information.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/layers.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/layout.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/link.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/music.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/package.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/package_go.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_code.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_excel.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_green.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_acrobat.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_c.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_cplusplus.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_flash.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_magnify.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_picture.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_powerpoint.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_stack.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_text.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_word.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_white_zip.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/page_world.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/photo.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/picture.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/plugin.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/plugin_add.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/script.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/script_gear.png
-rw-r--r--	root/root	usr/share/squid/icons/silk/script_palette.png
-rw-r--r--	root/root	usr/share/squid/mib.txt
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwx------	squid/squid	var/log/squid/
drwxr-xr-x	root/root	var/run/
drwx------	squid/squid	var/run/squid/
drwx------	squid/squid	var/squid/
