# ports/contrib/thttpd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/thttpd
-rw-r--r--	root/root	etc/thttpd.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/htpasswd
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/syslogtocern
-rwxr-xr-x	root/root	usr/sbin/thttpd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rwxr-xr-x	root/root	usr/share/man/man1/htpasswd.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rwxr-xr-x	root/root	usr/share/man/man8/redirect.8.gz
-rwxr-xr-x	root/root	usr/share/man/man8/ssi.8.gz
-rwxr-xr-x	root/root	usr/share/man/man8/syslogtocern.8.gz
-rw-r--r--	root/root	usr/share/man/man8/thttpd.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
-rw-r--r--	root/root	var/log/thttpd.log (EMPTY)
drwxr-xr-x	root/root	var/www/
drwxr-xr-x	root/root	var/www/cgi-bin/
-rwxr-xr-x	root/root	var/www/cgi-bin/phf
-rwxr-xr-x	root/root	var/www/cgi-bin/printenv
-rwxr-xr-x	root/root	var/www/cgi-bin/redirect
-rwxr-xr-x	root/root	var/www/cgi-bin/ssi
drwxrwxrwt	root/root	var/www/users/
