# image/ports/xorg/xorg-xdm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/X11/
drwxr-xr-x	root/root	etc/X11/xdm/
-rw-r--r--	root/root	etc/X11/xdm/Xaccess
-rw-r--r--	root/root	etc/X11/xdm/Xresources
-rw-r--r--	root/root	etc/X11/xdm/Xservers
-rw-r--r--	root/root	etc/X11/xdm/xdm-config
drwxr-xr-x	root/root	etc/pam.d/
-rwxr-xr-x	root/root	etc/pam.d/xdm
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/xdm
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/xdm
-rwxr-xr-x	root/root	usr/bin/xdmshell
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/X11/
drwxr-xr-x	root/root	usr/lib/X11/xdm/
-rwxr-xr-x	root/root	usr/lib/X11/xdm/GiveConsole
-rwxr-xr-x	root/root	usr/lib/X11/xdm/TakeConsole
-rwxr-xr-x	root/root	usr/lib/X11/xdm/Xreset
-rwxr-xr-x	root/root	usr/lib/X11/xdm/Xsession
-rwxr-xr-x	root/root	usr/lib/X11/xdm/Xsetup_0
-rwxr-xr-x	root/root	usr/lib/X11/xdm/Xstartup
-rwxr-xr-x	root/root	usr/lib/X11/xdm/Xwilling
-rwxr-xr-x	root/root	usr/lib/X11/xdm/chooser
-rwxr-xr-x	root/root	usr/lib/X11/xdm/libXdmGreet.la
-rwxr-xr-x	root/root	usr/lib/X11/xdm/libXdmGreet.so
drwxr-xr-x	root/root	usr/lib/X11/xdm/pixmaps/
-rw-r--r--	root/root	usr/lib/X11/xdm/pixmaps/xorg-bw.xpm
-rw-r--r--	root/root	usr/lib/X11/xdm/pixmaps/xorg.xpm
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/X11/
drwxr-xr-x	root/root	usr/share/X11/app-defaults/
-rw-r--r--	root/root	usr/share/X11/app-defaults/Chooser
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/xdm.8.gz
-rw-r--r--	root/root	usr/share/man/man8/xdmshell.8.gz
