# image/ports/opt/bitlbee/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/bitlbee/
-rw-r--r--	root/root	etc/bitlbee/bitlbee.conf
-rw-r--r--	root/root	etc/bitlbee/motd.txt
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/bitlbee
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/bitlbee
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/bitlbee/
-rw-r--r--	root/root	usr/share/bitlbee/help.txt
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/bitlbee.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/bitlbee.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	bitlbee/root	var/lib/bitlbee/
