# image/ports/contrib/openvswitch/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/bash_completion.d/
-rwxr-xr-x	root/root	etc/bash_completion.d/ovs-appctl-bashcomp.bash
-rwxr-xr-x	root/root	etc/bash_completion.d/ovs-vsctl-bashcomp.bash
drwxr-xr-x	root/root	etc/openvswitch/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/ovsd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ovs-appctl
-rwxr-xr-x	root/root	usr/bin/ovs-docker
-rwxr-xr-x	root/root	usr/bin/ovs-dpctl
-rwxr-xr-x	root/root	usr/bin/ovs-dpctl-top
-rwxr-xr-x	root/root	usr/bin/ovs-l3ping
-rwxr-xr-x	root/root	usr/bin/ovs-ofctl
-rwxr-xr-x	root/root	usr/bin/ovs-parse-backtrace
-rwxr-xr-x	root/root	usr/bin/ovs-pcap
-rwxr-xr-x	root/root	usr/bin/ovs-pki
-rwxr-xr-x	root/root	usr/bin/ovs-tcpdump
-rwxr-xr-x	root/root	usr/bin/ovs-tcpundump
-rwxr-xr-x	root/root	usr/bin/ovs-test
-rwxr-xr-x	root/root	usr/bin/ovs-testcontroller
-rwxr-xr-x	root/root	usr/bin/ovs-vlan-test
-rwxr-xr-x	root/root	usr/bin/ovs-vsctl
-rwxr-xr-x	root/root	usr/bin/ovsdb-client
-rwxr-xr-x	root/root	usr/bin/ovsdb-tool
-rwxr-xr-x	root/root	usr/bin/vtep-ctl
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/openflow/
-rw-r--r--	root/root	usr/include/openflow/intel-ext.h
-rw-r--r--	root/root	usr/include/openflow/netronome-ext.h
-rw-r--r--	root/root	usr/include/openflow/nicira-ext.h
-rw-r--r--	root/root	usr/include/openflow/openflow-1.0.h
-rw-r--r--	root/root	usr/include/openflow/openflow-1.1.h
-rw-r--r--	root/root	usr/include/openflow/openflow-1.2.h
-rw-r--r--	root/root	usr/include/openflow/openflow-1.3.h
-rw-r--r--	root/root	usr/include/openflow/openflow-1.4.h
-rw-r--r--	root/root	usr/include/openflow/openflow-1.5.h
-rw-r--r--	root/root	usr/include/openflow/openflow-common.h
-rw-r--r--	root/root	usr/include/openflow/openflow.h
drwxr-xr-x	root/root	usr/include/openvswitch/
-rw-r--r--	root/root	usr/include/openvswitch/compiler.h
-rw-r--r--	root/root	usr/include/openvswitch/dynamic-string.h
-rw-r--r--	root/root	usr/include/openvswitch/flow.h
-rw-r--r--	root/root	usr/include/openvswitch/geneve.h
-rw-r--r--	root/root	usr/include/openvswitch/hmap.h
-rw-r--r--	root/root	usr/include/openvswitch/json.h
-rw-r--r--	root/root	usr/include/openvswitch/list.h
-rw-r--r--	root/root	usr/include/openvswitch/match.h
-rw-r--r--	root/root	usr/include/openvswitch/meta-flow.h
-rw-r--r--	root/root	usr/include/openvswitch/namemap.h
-rw-r--r--	root/root	usr/include/openvswitch/netdev.h
-rw-r--r--	root/root	usr/include/openvswitch/nsh.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-actions.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-bundle.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-connection.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-ed-props.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-errors.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-flow.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-group.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-ipfix.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-match.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-meter.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-monitor.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-msgs.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-packet.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-parse.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-port.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-print.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-prop.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-protocol.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-queue.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-switch.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-table.h
-rw-r--r--	root/root	usr/include/openvswitch/ofp-util.h
-rw-r--r--	root/root	usr/include/openvswitch/ofpbuf.h
-rw-r--r--	root/root	usr/include/openvswitch/packets.h
-rw-r--r--	root/root	usr/include/openvswitch/poll-loop.h
-rw-r--r--	root/root	usr/include/openvswitch/rconn.h
-rw-r--r--	root/root	usr/include/openvswitch/shash.h
-rw-r--r--	root/root	usr/include/openvswitch/thread.h
-rw-r--r--	root/root	usr/include/openvswitch/token-bucket.h
-rw-r--r--	root/root	usr/include/openvswitch/tun-metadata.h
-rw-r--r--	root/root	usr/include/openvswitch/type-props.h
-rw-r--r--	root/root	usr/include/openvswitch/types.h
-rw-r--r--	root/root	usr/include/openvswitch/util.h
-rw-r--r--	root/root	usr/include/openvswitch/uuid.h
-rw-r--r--	root/root	usr/include/openvswitch/vconn.h
-rw-r--r--	root/root	usr/include/openvswitch/version.h
-rw-r--r--	root/root	usr/include/openvswitch/vlog.h
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libofproto.a
-rwxr-xr-x	root/root	usr/lib/libofproto.la
-rw-r--r--	root/root	usr/lib/libopenvswitch.a
-rwxr-xr-x	root/root	usr/lib/libopenvswitch.la
-rw-r--r--	root/root	usr/lib/libopenvswitchavx512.a
-rwxr-xr-x	root/root	usr/lib/libopenvswitchavx512.la
-rw-r--r--	root/root	usr/lib/libovsdb.a
-rwxr-xr-x	root/root	usr/lib/libovsdb.la
-rw-r--r--	root/root	usr/lib/libsflow.a
-rwxr-xr-x	root/root	usr/lib/libsflow.la
-rw-r--r--	root/root	usr/lib/libvtep.a
-rwxr-xr-x	root/root	usr/lib/libvtep.la
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libofproto.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libopenvswitch.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libovsdb.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libsflow.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/ovs-bugtool
-rwxr-xr-x	root/root	usr/sbin/ovs-vswitchd
-rwxr-xr-x	root/root	usr/sbin/ovsdb-server
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/ovs-pcap.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ovsdb-client.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ovsdb-server.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ovsdb-tool.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/ovs-vswitchd.conf.db.5.gz
-rw-r--r--	root/root	usr/share/man/man5/ovsdb-server.5.gz
-rw-r--r--	root/root	usr/share/man/man5/vtep.5.gz
drwxr-xr-x	root/root	usr/share/man/man7/
-rw-r--r--	root/root	usr/share/man/man7/ovs-actions.7.gz
-rw-r--r--	root/root	usr/share/man/man7/ovs-fields.7.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/ovs-bugtool.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ovs-dpctl-top.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ovs-dpctl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ovs-kmod-ctl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ovs-ofctl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ovs-testcontroller.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ovs-vsctl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ovs-vswitchd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/vtep-ctl.8.gz
drwxr-xr-x	root/root	usr/share/openvswitch/
drwxr-xr-x	root/root	usr/share/openvswitch/bugtool-plugins/
drwxr-xr-x	root/root	usr/share/openvswitch/bugtool-plugins/kernel-info/
-rw-r--r--	root/root	usr/share/openvswitch/bugtool-plugins/kernel-info/openvswitch.xml
drwxr-xr-x	root/root	usr/share/openvswitch/bugtool-plugins/network-status/
-rw-r--r--	root/root	usr/share/openvswitch/bugtool-plugins/network-status/openvswitch.xml
-rw-r--r--	root/root	usr/share/openvswitch/bugtool-plugins/system-configuration.xml
drwxr-xr-x	root/root	usr/share/openvswitch/bugtool-plugins/system-configuration/
-rw-r--r--	root/root	usr/share/openvswitch/bugtool-plugins/system-configuration/openvswitch.xml
drwxr-xr-x	root/root	usr/share/openvswitch/bugtool-plugins/system-logs/
-rw-r--r--	root/root	usr/share/openvswitch/bugtool-plugins/system-logs/openvswitch.xml
drwxr-xr-x	root/root	usr/share/openvswitch/python/
drwxr-xr-x	root/root	usr/share/openvswitch/python/ovs/
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/__init__.py
drwxr-xr-x	root/root	usr/share/openvswitch/python/ovs/compat/
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/compat/__init__.py (EMPTY)
drwxr-xr-x	root/root	usr/share/openvswitch/python/ovs/compat/sortedcontainers/
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/compat/sortedcontainers/__init__.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/compat/sortedcontainers/sorteddict.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/compat/sortedcontainers/sortedlist.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/compat/sortedcontainers/sortedset.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/daemon.py
drwxr-xr-x	root/root	usr/share/openvswitch/python/ovs/db/
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/__init__.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/custom_index.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/data.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/error.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/idl.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/parser.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/schema.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/db/types.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/dirs.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/fatal_signal.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/fcntl_win.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/json.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/jsonrpc.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/ovsuuid.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/poller.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/process.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/reconnect.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/socket_util.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/stream.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/timeval.py
drwxr-xr-x	root/root	usr/share/openvswitch/python/ovs/unixctl/
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/unixctl/__init__.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/unixctl/client.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/unixctl/server.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/util.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/version.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/vlog.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovs/winutils.py
drwxr-xr-x	root/root	usr/share/openvswitch/python/ovstest/
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/__init__.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/args.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/rpcserver.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/tcp.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/tests.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/udp.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/util.py
-rw-r--r--	root/root	usr/share/openvswitch/python/ovstest/vswitch.py
drwxr-xr-x	root/root	usr/share/openvswitch/scripts/
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-daemons-ver
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-fdb-show
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-get-dpdk-nic-numa
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-get-port-stats
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-ovs-appctl-dpif
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-ovs-bridge-datapath-type
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-ovs-ofctl-loop-over-bridges
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-ovs-vswitchd-threads-affinity
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-qos-configs
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-bugtool-tc-class-show
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-check-dead-ifs
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-ctl
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-kmod-ctl
-rw-r--r--	root/root	usr/share/openvswitch/scripts/ovs-lib
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-monitor-ipsec
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-save
-rwxr-xr-x	root/root	usr/share/openvswitch/scripts/ovs-vtep
-rw-r--r--	root/root	usr/share/openvswitch/vswitch.ovsschema
-rw-r--r--	root/root	usr/share/openvswitch/vtep.ovsschema
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/openvswitch/
drwxr-xr-x	root/root	var/lib/openvswitch/pki/
