# image/ports/opt/rsyslog/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/rsyslogd
-rw-r--r--	root/root	etc/rsyslog.conf
drwxr-xr-x	root/root	etc/rsyslog.d/
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/rsyslog/
-rwxr-xr-x	root/root	usr/lib/rsyslog/fmhash.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/fmhash.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/fmhttp.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/fmhttp.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/imklog.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/imklog.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/immark.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/immark.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/imtcp.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/imtcp.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/imudp.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/imudp.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/imuxsock.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/imuxsock.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmcry_gcry.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmcry_gcry.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmnet.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmnet.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmnetstrms.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmnetstrms.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmnsd_ptcp.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmnsd_ptcp.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmregexp.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmregexp.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmtcpclt.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmtcpclt.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmtcpsrv.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmtcpsrv.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmzlibw.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/lmzlibw.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/mmexternal.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/mmexternal.so
-rwxr-xr-x	root/root	usr/lib/rsyslog/omtesting.la
-rwxr-xr-x	root/root	usr/lib/rsyslog/omtesting.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/rsyslogd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/rsyslog.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/rsyslogd.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	root/root	var/lib/rsyslog/
