# ports/opt/unbound/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/unbound
drwxr-xr-x	root/root	etc/unbound/
drwxr-xr-x	unbound/unbound	etc/unbound/anchor/
-rw-r--r--	root/root	etc/unbound/unbound.conf
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/unbound
-rwxr-xr-x	root/root	usr/sbin/unbound-anchor
-rwxr-xr-x	root/root	usr/sbin/unbound-checkconf
-rwxr-xr-x	root/root	usr/sbin/unbound-control
-rwxr-xr-x	root/root	usr/sbin/unbound-control-setup
-rwxr-xr-x	root/root	usr/sbin/unbound-host
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/unbound-host.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/unbound.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/unbound-anchor.8.gz
-rw-r--r--	root/root	usr/share/man/man8/unbound-checkconf.8.gz
-rw-r--r--	root/root	usr/share/man/man8/unbound-control-setup.8.gz
-rw-r--r--	root/root	usr/share/man/man8/unbound-control.8.gz
-rw-r--r--	root/root	usr/share/man/man8/unbound.8.gz
