# ports/opt/postfix/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/postfix/
-rw-r--r--	root/root	etc/postfix/access
-rw-r--r--	root/root	etc/postfix/access.db (EMPTY)
-rw-r--r--	root/root	etc/postfix/aliases
-rw-r--r--	root/root	etc/postfix/aliases.db (EMPTY)
-rw-r--r--	root/root	etc/postfix/body_checks (EMPTY)
-rw-r--r--	root/root	etc/postfix/bounce.cf.default
-rw-r--r--	root/root	etc/postfix/canonical
-rw-r--r--	root/root	etc/postfix/canonical.db (EMPTY)
-rw-r--r--	root/root	etc/postfix/generic
-rw-r--r--	root/root	etc/postfix/header_checks
-rw-r--r--	root/root	etc/postfix/main.cf
-rw-r--r--	root/root	etc/postfix/main.cf.default
-rw-r--r--	root/root	etc/postfix/main.cf.proto
-rw-r--r--	root/root	etc/postfix/makedefs.out
-rw-r--r--	root/root	etc/postfix/master.cf
-rw-r--r--	root/root	etc/postfix/master.cf.proto
-rw-r--r--	root/root	etc/postfix/mime_header_checks (EMPTY)
-rw-r--r--	root/root	etc/postfix/postfix-files
drwxr-xr-x	root/root	etc/postfix/postfix-files.d/
-rw-r--r--	root/root	etc/postfix/relay_domains (EMPTY)
-rw-r--r--	root/root	etc/postfix/relocated
-rw-r--r--	root/root	etc/postfix/relocated.db (EMPTY)
-rw-r--r--	root/root	etc/postfix/transport
-rw-r--r--	root/root	etc/postfix/transport.db (EMPTY)
-rw-r--r--	root/root	etc/postfix/virtual
-rw-r--r--	root/root	etc/postfix/virtual.db (EMPTY)
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/postfix
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
lrwxrwxrwx	root/root	usr/bin/mailq -> ../../usr/sbin/sendmail
lrwxrwxrwx	root/root	usr/bin/newaliases -> ../../usr/sbin/sendmail
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/postfix/
-rwxr-xr-x	root/root	usr/lib/postfix/anvil
-rwxr-xr-x	root/root	usr/lib/postfix/bounce
-rwxr-xr-x	root/root	usr/lib/postfix/cleanup
-rwxr-xr-x	root/root	usr/lib/postfix/discard
-rwxr-xr-x	root/root	usr/lib/postfix/dnsblog
-rwxr-xr-x	root/root	usr/lib/postfix/error
-rwxr-xr-x	root/root	usr/lib/postfix/flush
-rwxr-xr-x	root/root	usr/lib/postfix/lmtp
-rwxr-xr-x	root/root	usr/lib/postfix/local
-rwxr-xr-x	root/root	usr/lib/postfix/master
-rwxr-xr-x	root/root	usr/lib/postfix/nqmgr
-rwxr-xr-x	root/root	usr/lib/postfix/oqmgr
-rwxr-xr-x	root/root	usr/lib/postfix/pickup
-rwxr-xr-x	root/root	usr/lib/postfix/pipe
-rwxr-xr-x	root/root	usr/lib/postfix/post-install
-rwxr-xr-x	root/root	usr/lib/postfix/postfix-script
-rwxr-xr-x	root/root	usr/lib/postfix/postfix-tls-script
-rwxr-xr-x	root/root	usr/lib/postfix/postfix-wrapper
-rwxr-xr-x	root/root	usr/lib/postfix/postlogd
-rwxr-xr-x	root/root	usr/lib/postfix/postmulti-script
-rwxr-xr-x	root/root	usr/lib/postfix/postscreen
-rwxr-xr-x	root/root	usr/lib/postfix/proxymap
-rwxr-xr-x	root/root	usr/lib/postfix/qmgr
-rwxr-xr-x	root/root	usr/lib/postfix/qmqpd
-rwxr-xr-x	root/root	usr/lib/postfix/scache
-rwxr-xr-x	root/root	usr/lib/postfix/showq
-rwxr-xr-x	root/root	usr/lib/postfix/smtp
-rwxr-xr-x	root/root	usr/lib/postfix/smtpd
-rwxr-xr-x	root/root	usr/lib/postfix/spawn
-rwxr-xr-x	root/root	usr/lib/postfix/tlsmgr
-rwxr-xr-x	root/root	usr/lib/postfix/tlsproxy
-rwxr-xr-x	root/root	usr/lib/postfix/trivial-rewrite
-rwxr-xr-x	root/root	usr/lib/postfix/verify
-rwxr-xr-x	root/root	usr/lib/postfix/virtual
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/postalias
-rwxr-xr-x	root/root	usr/sbin/postcat
-rwxr-xr-x	root/root	usr/sbin/postconf
-rwxr-sr-x	root/postdrop	usr/sbin/postdrop
-rwxr-xr-x	root/root	usr/sbin/postfix
-rwxr-xr-x	root/root	usr/sbin/postkick
-rwxr-xr-x	root/root	usr/sbin/postlock
-rwxr-xr-x	root/root	usr/sbin/postlog
-rwxr-xr-x	root/root	usr/sbin/postmap
-rwxr-xr-x	root/root	usr/sbin/postmulti
-rwxr-sr-x	root/postdrop	usr/sbin/postqueue
-rwxr-xr-x	root/root	usr/sbin/postsuper
-rwxr-xr-x	root/root	usr/sbin/sendmail
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/mailq.1.gz
-rw-r--r--	root/root	usr/share/man/man1/newaliases.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postalias.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postcat.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postconf.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postdrop.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postfix-tls.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postfix.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postkick.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postlock.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postlog.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postmap.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postmulti.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postqueue.1.gz
-rw-r--r--	root/root	usr/share/man/man1/postsuper.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sendmail.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/access.5.gz
-rw-r--r--	root/root	usr/share/man/man5/aliases.5.gz
-rw-r--r--	root/root	usr/share/man/man5/body_checks.5.gz
-rw-r--r--	root/root	usr/share/man/man5/bounce.5.gz
-rw-r--r--	root/root	usr/share/man/man5/canonical.5.gz
-rw-r--r--	root/root	usr/share/man/man5/cidr_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/generic.5.gz
-rw-r--r--	root/root	usr/share/man/man5/header_checks.5.gz
-rw-r--r--	root/root	usr/share/man/man5/ldap_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/lmdb_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/master.5.gz
-rw-r--r--	root/root	usr/share/man/man5/memcache_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/mysql_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/nisplus_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/pcre_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/pgsql_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/postconf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/postfix-wrapper.5.gz
-rw-r--r--	root/root	usr/share/man/man5/regexp_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/relocated.5.gz
-rw-r--r--	root/root	usr/share/man/man5/socketmap_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/sqlite_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/tcp_table.5.gz
-rw-r--r--	root/root	usr/share/man/man5/transport.5.gz
-rw-r--r--	root/root	usr/share/man/man5/virtual.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/anvil.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bounce.8.gz
-rw-r--r--	root/root	usr/share/man/man8/cleanup.8.gz
-rw-r--r--	root/root	usr/share/man/man8/defer.8.gz
-rw-r--r--	root/root	usr/share/man/man8/discard.8.gz
-rw-r--r--	root/root	usr/share/man/man8/dnsblog.8.gz
-rw-r--r--	root/root	usr/share/man/man8/error.8.gz
-rw-r--r--	root/root	usr/share/man/man8/flush.8.gz
-rw-r--r--	root/root	usr/share/man/man8/lmtp.8.gz
-rw-r--r--	root/root	usr/share/man/man8/local.8.gz
-rw-r--r--	root/root	usr/share/man/man8/master.8.gz
-rw-r--r--	root/root	usr/share/man/man8/oqmgr.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pickup.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pipe.8.gz
-rw-r--r--	root/root	usr/share/man/man8/postlogd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/postscreen.8.gz
-rw-r--r--	root/root	usr/share/man/man8/proxymap.8.gz
-rw-r--r--	root/root	usr/share/man/man8/qmgr.8.gz
-rw-r--r--	root/root	usr/share/man/man8/qmqpd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/scache.8.gz
-rw-r--r--	root/root	usr/share/man/man8/showq.8.gz
-rw-r--r--	root/root	usr/share/man/man8/smtp.8.gz
-rw-r--r--	root/root	usr/share/man/man8/smtpd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/spawn.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tlsmgr.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tlsproxy.8.gz
-rw-r--r--	root/root	usr/share/man/man8/trace.8.gz
-rw-r--r--	root/root	usr/share/man/man8/trivial-rewrite.8.gz
-rw-r--r--	root/root	usr/share/man/man8/verify.8.gz
-rw-r--r--	root/root	usr/share/man/man8/virtual.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwx------	mail/root	var/lib/postfix/
drwxr-xr-x	root/root	var/spool/
drwxr-xr-x	root/root	var/spool/postfix/
drwx------	mail/root	var/spool/postfix/active/
drwx------	mail/root	var/spool/postfix/bounce/
drwx------	mail/root	var/spool/postfix/corrupt/
drwx------	mail/root	var/spool/postfix/defer/
drwx------	mail/root	var/spool/postfix/deferred/
drwx------	mail/root	var/spool/postfix/flush/
drwx------	mail/root	var/spool/postfix/hold/
drwx------	mail/root	var/spool/postfix/incoming/
drwx-wx---	mail/postdrop	var/spool/postfix/maildrop/
drwxr-xr-x	mail/root	var/spool/postfix/pid/
drwx------	mail/root	var/spool/postfix/private/
drwx--x---	mail/postdrop	var/spool/postfix/public/
drwx------	mail/root	var/spool/postfix/saved/
drwx------	mail/root	var/spool/postfix/trace/
