# image/ports/contrib/bumblebee/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/bash_completion.d/
-rw-r--r--	root/root	etc/bash_completion.d/optirun
drwxr-xr-x	root/root	etc/bumblebee/
-rw-r--r--	root/root	etc/bumblebee/bumblebee.conf
drwxr-xr-x	root/root	etc/bumblebee/xorg.conf.d/
-rw-r--r--	root/root	etc/bumblebee/xorg.conf.d/10-dummy.conf (EMPTY)
-rw-r--r--	root/root	etc/bumblebee/xorg.conf.nouveau
-rw-r--r--	root/root	etc/bumblebee/xorg.conf.nvidia
drwxr-xr-x	root/root	etc/modprobe.d/
-rw-r--r--	root/root	etc/modprobe.d/bumblebee.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/bumblebeed
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/99-bumblebee-nvidia-dev.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/bumblebee-bugreport
-rwxr-xr-x	root/root	usr/bin/optirun
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/bumblebeed
drwxr-xr-x	root/root	usr/share/
