# image/ports/contrib/hefur/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/hefurd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/hefurctl
-rwxr-xr-x	root/root	usr/bin/hefurd
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/hefur/
drwxr-xr-x	root/root	usr/share/hefur/torrents/
drwxr-xr-x	root/root	usr/share/hefur/www/
drwxr-xr-x	root/root	usr/share/hefur/www/css/
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-grid.css
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-grid.css.map
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-grid.min.css
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-grid.min.css.map
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-reboot.css
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-reboot.css.map
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-reboot.min.css
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap-reboot.min.css.map
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap.css
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap.css.map
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap.min.css
-rw-r--r--	root/root	usr/share/hefur/www/css/bootstrap.min.css.map
drwxr-xr-x	root/root	usr/share/hefur/www/img/
-rw-r--r--	root/root	usr/share/hefur/www/img/favicon.ico
-rw-r--r--	root/root	usr/share/hefur/www/img/magnet-icon.gif
-rw-r--r--	root/root	usr/share/hefur/www/img/peers.png
-rw-r--r--	root/root	usr/share/hefur/www/img/torrent-icon.png
drwxr-xr-x	root/root	usr/share/hefur/www/js/
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.bundle.js
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.bundle.js.map
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.bundle.min.js
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.bundle.min.js.map
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.js
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.js.map
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.min.js
-rw-r--r--	root/root	usr/share/hefur/www/js/bootstrap.min.js.map
-rw-r--r--	root/root	usr/share/hefur/www/js/table-sorter.js
drwxr-xr-x	root/root	usr/share/hefur/www/tpl/
-rw-r--r--	root/root	usr/share/hefur/www/tpl/page.html
-rw-r--r--	root/root	usr/share/hefur/www/tpl/peers.html
-rw-r--r--	root/root	usr/share/hefur/www/tpl/stat.html
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
-rw-r--r--	root/root	var/log/hefurd.log (EMPTY)
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	root/root	var/run/hefur/
