# image/ports/contrib/nbd/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/nbd-server/
-rw-------	root/root	etc/nbd-server/allow
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/nbd-server
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/nbd-server
-rwxr-xr-x	root/root	usr/bin/nbd-trdump
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/nbd-client
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/nbd-server.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nbd-trdump.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/nbd-server.5.gz
-rw-r--r--	root/root	usr/share/man/man5/nbdtab.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/nbd-client.8.gz
