# image/ports/contrib/nut/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/nut/
-rw-r--r--	root/nut	etc/nut/nut.conf.sample
-rw-r--r--	root/nut	etc/nut/ups.conf.sample
-rw-r--r--	root/nut	etc/nut/upsd.conf.sample
-rw-r--r--	root/nut	etc/nut/upsd.users.sample
-rw-r--r--	root/nut	etc/nut/upsmon.conf.sample
-rw-r--r--	root/nut	etc/nut/upssched.conf.sample
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr--	root/root	etc/rc.d/upsd
drwxr-xr-x	root/root	etc/udev/
drwxr-xr-x	root/root	etc/udev/rules.d/
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/al175
-rwxr-xr-x	root/root	usr/bin/apcsmart
-rwxr-xr-x	root/root	usr/bin/apcsmart-old
-rwxr-xr-x	root/root	usr/bin/apcupsd-ups
-rwxr-xr-x	root/root	usr/bin/bcmxcp
-rwxr-xr-x	root/root	usr/bin/belkin
-rwxr-xr-x	root/root	usr/bin/belkinunv
-rwxr-xr-x	root/root	usr/bin/bestfcom
-rwxr-xr-x	root/root	usr/bin/bestfortress
-rwxr-xr-x	root/root	usr/bin/bestuferrups
-rwxr-xr-x	root/root	usr/bin/bestups
-rwxr-xr-x	root/root	usr/bin/blazer_ser
-rwxr-xr-x	root/root	usr/bin/clone
-rwxr-xr-x	root/root	usr/bin/clone-outlet
-rwxr-xr-x	root/root	usr/bin/dummy-ups
-rwxr-xr-x	root/root	usr/bin/etapro
-rwxr-xr-x	root/root	usr/bin/everups
-rwxr-xr-x	root/root	usr/bin/gamatronic
-rwxr-xr-x	root/root	usr/bin/genericups
-rwxr-xr-x	root/root	usr/bin/isbmex
-rwxr-xr-x	root/root	usr/bin/ivtscd
-rwxr-xr-x	root/root	usr/bin/liebert
-rwxr-xr-x	root/root	usr/bin/liebert-esp2
-rwxr-xr-x	root/root	usr/bin/masterguard
-rwxr-xr-x	root/root	usr/bin/metasys
-rwxr-xr-x	root/root	usr/bin/mge-shut
-rwxr-xr-x	root/root	usr/bin/mge-utalk
-rwxr-xr-x	root/root	usr/bin/microdowell
-rwxr-xr-x	root/root	usr/bin/nut-scanner
-rwxr-xr-x	root/root	usr/bin/nutdrv_qx
-rwxr-xr-x	root/root	usr/bin/oldmge-shut
-rwxr-xr-x	root/root	usr/bin/oneac
-rwxr-xr-x	root/root	usr/bin/optiups
-rwxr-xr-x	root/root	usr/bin/powercom
-rwxr-xr-x	root/root	usr/bin/powerpanel
-rwxr-xr-x	root/root	usr/bin/rhino
-rwxr-xr-x	root/root	usr/bin/riello_ser
-rwxr-xr-x	root/root	usr/bin/safenet
-rwxr-xr-x	root/root	usr/bin/skel
-rwxr-xr-x	root/root	usr/bin/solis
-rwxr-xr-x	root/root	usr/bin/tripplite
-rwxr-xr-x	root/root	usr/bin/tripplitesu
-rwxr-xr-x	root/root	usr/bin/upsc
-rwxr-xr-x	root/root	usr/bin/upscmd
-rwxr-xr-x	root/root	usr/bin/upscode2
-rwxr-xr-x	root/root	usr/bin/upslog
-rwxr-xr-x	root/root	usr/bin/upsrw
-rwxr-xr-x	root/root	usr/bin/upssched-cmd
-rwxr-xr-x	root/root	usr/bin/victronups
drwxr-xr-x	root/root	usr/lib/
-rw-r--r--	root/root	usr/lib/libnutclient.a
-rwxr-xr-x	root/root	usr/lib/libnutclient.la
lrwxrwxrwx	root/root	usr/lib/libnutclient.so -> libnutclient.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libnutclient.so.0 -> libnutclient.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libnutclient.so.0.0.0
-rw-r--r--	root/root	usr/lib/libnutscan.a
-rwxr-xr-x	root/root	usr/lib/libnutscan.la
lrwxrwxrwx	root/root	usr/lib/libnutscan.so -> libnutscan.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libnutscan.so.1 -> libnutscan.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libnutscan.so.1.0.0
-rw-r--r--	root/root	usr/lib/libupsclient.a
-rwxr-xr-x	root/root	usr/lib/libupsclient.la
lrwxrwxrwx	root/root	usr/lib/libupsclient.so -> libupsclient.so.4.0.0
lrwxrwxrwx	root/root	usr/lib/libupsclient.so.4 -> libupsclient.so.4.0.0
-rwxr-xr-x	root/root	usr/lib/libupsclient.so.4.0.0
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/upsd
-rwxr-xr-x	root/root	usr/sbin/upsdrvctl
-rwxr-xr-x	root/root	usr/sbin/upsmon
-rwxr-xr-x	root/root	usr/sbin/upssched
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/nut.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/ups.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/upsd.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/upsd.users.5.gz
-rw-r--r--	root/root	usr/share/man/man5/upsmon.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/upssched.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/al175.8.gz
-rw-r--r--	root/root	usr/share/man/man8/apcsmart-old.8.gz
-rw-r--r--	root/root	usr/share/man/man8/apcsmart.8.gz
-rw-r--r--	root/root	usr/share/man/man8/apcupsd-ups.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bcmxcp.8.gz
-rw-r--r--	root/root	usr/share/man/man8/belkin.8.gz
-rw-r--r--	root/root	usr/share/man/man8/belkinunv.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bestfcom.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bestfortress.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bestuferrups.8.gz
-rw-r--r--	root/root	usr/share/man/man8/bestups.8.gz
-rw-r--r--	root/root	usr/share/man/man8/blazer_ser.8.gz
-rw-r--r--	root/root	usr/share/man/man8/clone.8.gz
-rw-r--r--	root/root	usr/share/man/man8/dummy-ups.8.gz
-rw-r--r--	root/root	usr/share/man/man8/etapro.8.gz
-rw-r--r--	root/root	usr/share/man/man8/everups.8.gz
-rw-r--r--	root/root	usr/share/man/man8/gamatronic.8.gz
-rw-r--r--	root/root	usr/share/man/man8/genericups.8.gz
-rw-r--r--	root/root	usr/share/man/man8/isbmex.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ivtscd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/liebert-esp2.8.gz
-rw-r--r--	root/root	usr/share/man/man8/liebert.8.gz
-rw-r--r--	root/root	usr/share/man/man8/masterguard.8.gz
-rw-r--r--	root/root	usr/share/man/man8/metasys.8.gz
-rw-r--r--	root/root	usr/share/man/man8/mge-shut.8.gz
-rw-r--r--	root/root	usr/share/man/man8/mge-utalk.8.gz
-rw-r--r--	root/root	usr/share/man/man8/microdowell.8.gz
-rw-r--r--	root/root	usr/share/man/man8/nut-recorder.8.gz
-rw-r--r--	root/root	usr/share/man/man8/nut-scanner.8.gz
-rw-r--r--	root/root	usr/share/man/man8/nutdrv_qx.8.gz
-rw-r--r--	root/root	usr/share/man/man8/nutupsdrv.8.gz
-rw-r--r--	root/root	usr/share/man/man8/oneac.8.gz
-rw-r--r--	root/root	usr/share/man/man8/optiups.8.gz
-rw-r--r--	root/root	usr/share/man/man8/powercom.8.gz
-rw-r--r--	root/root	usr/share/man/man8/powerpanel.8.gz
-rw-r--r--	root/root	usr/share/man/man8/rhino.8.gz
-rw-r--r--	root/root	usr/share/man/man8/riello_ser.8.gz
-rw-r--r--	root/root	usr/share/man/man8/safenet.8.gz
-rw-r--r--	root/root	usr/share/man/man8/solis.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tripplite.8.gz
-rw-r--r--	root/root	usr/share/man/man8/tripplitesu.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upsc.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upscmd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upscode2.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upsd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upsdrvctl.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upslog.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upsmon.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upsrw.8.gz
-rw-r--r--	root/root	usr/share/man/man8/upssched.8.gz
-rw-r--r--	root/root	usr/share/man/man8/victronups.8.gz
drwxr-xr-x	root/root	usr/share/nut/
-rw-r--r--	root/root	usr/share/nut/cmdvartab
-rw-r--r--	root/root	usr/share/nut/driver.list
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/state/
drwxrwx---	root/nut	var/state/ups/
