# ports/opt/pure-ftpd/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/ftpusers
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/pure-ftpd
-rw-r--r--	root/root	etc/pure-ftpd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/pure-ftpd
drwxr-xr-x	root/root	etc/ssl/
drwxr-xr-x	root/root	etc/ssl/certs/
-rw-------	root/root	etc/ssl/certs/pure-ftpd.pem (EMPTY)
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/pure-pw
-rwxr-xr-x	root/root	usr/bin/pure-pwconvert
-rwxr-xr-x	root/root	usr/bin/pure-statsdecode
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/pure-authd
-rwxr-xr-x	root/root	usr/sbin/pure-certd
-rwxr-xr-x	root/root	usr/sbin/pure-ftpd
-rwxr-xr-x	root/root	usr/sbin/pure-ftpwho
-rwxr-xr-x	root/root	usr/sbin/pure-mrtginfo
-rwxr-xr-x	root/root	usr/sbin/pure-quotacheck
-rwxr-xr-x	root/root	usr/sbin/pure-uploadscript
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/pure-authd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-certd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-ftpd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-ftpwho.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-mrtginfo.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-pw.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-pwconvert.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-quotacheck.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-statsdecode.8.gz
-rw-r--r--	root/root	usr/share/man/man8/pure-uploadscript.8.gz
