# ports/contrib/clamav/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/clamd.conf
-rw-r--r--	root/root	etc/freshclam.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/clamd
-rwxr-xr-x	root/root	etc/rc.d/freshclam
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/clamav-config
-rwxr-xr-x	root/root	usr/bin/clambc
-rwxr-xr-x	root/root	usr/bin/clamconf
-rwxr-xr-x	root/root	usr/bin/clamdscan
-rwxr-xr-x	root/root	usr/bin/clamdtop
-rwxr-xr-x	root/root	usr/bin/clamscan
-rwxr-xr-x	root/root	usr/bin/freshclam
-rwxr-xr-x	root/root	usr/bin/sigtool
drwxr-xr-x	root/root	usr/include/
-rw-r--r--	root/root	usr/include/clamav-types.h
-rw-r--r--	root/root	usr/include/clamav-version.h
-rw-r--r--	root/root	usr/include/clamav.h
-rw-r--r--	root/root	usr/include/libfreshclam.h
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/libclamav.la
lrwxrwxrwx	root/root	usr/lib/libclamav.so -> libclamav.so.9.0.5
lrwxrwxrwx	root/root	usr/lib/libclamav.so.9 -> libclamav.so.9.0.5
-rwxr-xr-x	root/root	usr/lib/libclamav.so.9.0.5
-rwxr-xr-x	root/root	usr/lib/libclammspack.la
lrwxrwxrwx	root/root	usr/lib/libclammspack.so -> libclammspack.so.0.1.0
lrwxrwxrwx	root/root	usr/lib/libclammspack.so.0 -> libclammspack.so.0.1.0
-rwxr-xr-x	root/root	usr/lib/libclammspack.so.0.1.0
-rwxr-xr-x	root/root	usr/lib/libclamunrar.la
lrwxrwxrwx	root/root	usr/lib/libclamunrar.so -> libclamunrar.so.9.0.5
lrwxrwxrwx	root/root	usr/lib/libclamunrar.so.9 -> libclamunrar.so.9.0.5
-rwxr-xr-x	root/root	usr/lib/libclamunrar.so.9.0.5
-rwxr-xr-x	root/root	usr/lib/libclamunrar_iface.la
lrwxrwxrwx	root/root	usr/lib/libclamunrar_iface.so -> libclamunrar_iface.so.9.0.5
lrwxrwxrwx	root/root	usr/lib/libclamunrar_iface.so.9 -> libclamunrar_iface.so.9.0.5
-rwxr-xr-x	root/root	usr/lib/libclamunrar_iface.so.9.0.5
-rwxr-xr-x	root/root	usr/lib/libfreshclam.la
lrwxrwxrwx	root/root	usr/lib/libfreshclam.so -> libfreshclam.so.2.0.1
lrwxrwxrwx	root/root	usr/lib/libfreshclam.so.2 -> libfreshclam.so.2.0.1
-rwxr-xr-x	root/root	usr/lib/libfreshclam.so.2.0.1
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/libclamav.pc
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/clamd
-rwxr-xr-x	root/root	usr/sbin/clamonacc
drwxr-xr-x	root/root	usr/share/
drwxrwxr-x	clamav/clamav	usr/share/clamav/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/clambc.1.gz
-rw-r--r--	root/root	usr/share/man/man1/clamconf.1.gz
-rw-r--r--	root/root	usr/share/man/man1/clamdscan.1.gz
-rw-r--r--	root/root	usr/share/man/man1/clamdtop.1.gz
-rw-r--r--	root/root	usr/share/man/man1/clamscan.1.gz
-rw-r--r--	root/root	usr/share/man/man1/freshclam.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sigtool.1.gz
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/clamav-milter.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/clamd.conf.5.gz
-rw-r--r--	root/root	usr/share/man/man5/freshclam.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/clamav-milter.8.gz
-rw-r--r--	root/root	usr/share/man/man8/clamd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/clamonacc.8.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	clamav/clamav	var/log/clamav/
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	clamav/clamav	var/run/clamav/
