# image/ports/contrib/conntrack-tools/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/conntrackd/
-rw-r--r--	root/root	etc/conntrackd/conntrackd.conf
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/conntrackd
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/conntrack-tools/
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_amanda.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_amanda.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_dhcpv6.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_dhcpv6.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_ftp.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_ftp.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_mdns.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_mdns.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_rpc.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_rpc.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_sane.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_sane.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_slp.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_slp.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_ssdp.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_ssdp.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_tftp.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_tftp.so
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_tns.la
-rwxr-xr-x	root/root	usr/lib/conntrack-tools/ct_helper_tns.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/conntrack
-rwxr-xr-x	root/root	usr/sbin/conntrackd
-rwxr-xr-x	root/root	usr/sbin/nfct
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man5/
-rw-r--r--	root/root	usr/share/man/man5/conntrackd.conf.5.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/conntrack.8.gz
-rw-r--r--	root/root	usr/share/man/man8/conntrackd.8.gz
-rw-r--r--	root/root	usr/share/man/man8/nfct.8.gz
